# aifmt.service.in
[Unit]
Description=AI format (binfmt_misc) support
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@pkgdatadir@/aifmt start
ExecStop=@pkgdatadir@/aifmt stop

[Install]
WantedBy=multi-user.target
